# packages/debian/vmpower-agent/etc/systemd/system/vmpower.service
[Unit]
Description=VMPower agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/vmpower 
User=vmpower

[Install]
WantedBy=multi-user.target
